# Raspberry/temperatures.supervisor.conf
[program:temperatures]
command=/usr/bin/python /home/pi/temps/main.py
directory=/home/pi/temps
redirect_stderr=true
startsecs=5
autorestart=true
stdout_logfile=/home/pi/temps/main.stdout.log
